# services/systemd-user/perf-api.service
[Unit]
Description=Frontend Performance Telemetry API
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /home/nikefd/perf-api.js
Restart=always
RestartSec=5
WorkingDirectory=/home/nikefd

[Install]
WantedBy=default.target
