# services/systemd-user/usage-api.service
[Unit]
Description=Token Usage API
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /home/nikefd/openclaw-deploy/services/usage/server.js
Restart=always
RestartSec=5
WorkingDirectory=/home/nikefd

[Install]
WantedBy=default.target
